# rootfs/profiles/sway/etc/systemd/user/mako.service
[Unit]
Description=Mako notification daemon
PartOf=sway-session.target graphical-session.target
After=sway-session.target

[Service]
ExecStart=/usr/bin/mako
Restart=on-failure

[Install]
WantedBy=sway-session.target
